# etc/udev/rules.d/81-thinkpad-dock.rules
SUBSYSTEM=="usb", ACTION=="add|remove", ENV{ID_VENDOR_ID}=="17ef", ENV{ID_MODEL_ID}=="1012", RUN+="/home/nalipaz/bin/thinkpad-dock"
